# infrastructure/roles/portal/templates/pnpu-portal.service.j2
[Unit]
Description=PNPU Portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ pnpu_service_user }}
Group={{ pnpu_service_group }}
WorkingDirectory={{ pnpu_base_dir }}/current
EnvironmentFile={{ pnpu_env_file }}
ExecStart={{ pnpu_node_bin }} run start
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths={{ pnpu_base_dir }}/shared
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
